# supervisor/fittrackee-workers.supervisor.conf
[program:fittrackee-workers]
command = .venv/bin/flask worker --processes %(ENV_WORKERS_PROCESSES)s

autostart = false
autorestart = false
